# system/etc/sv/lxc-containers/finish
#!/system/bin/sh
$PREFIX/libexec/lxc/lxc-containers stop
