# system/etc/sv/dockerd/run
#!/system/bin/sh
"PATH=$PATH /data/data/com.termux/files/usr/bin/dockerd 2>&1"